# packaging/launchd/com.github.chagui.besreceiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.chagui.besreceiver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/besreceiver</string>
        <string>--config=/usr/local/etc/besreceiver/config.yaml</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BESRECEIVER_CONFIG</key>
        <string>/usr/local/etc/besreceiver/config.yaml</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/besreceiver/besreceiver.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/besreceiver/besreceiver.err.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/besreceiver</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
</dict>
</plist>
